# rootfs/usr/lib/systemd/system/systemd-homed.service.d/10-firstboot.conf
[Unit]
Wants=systemd-homed-firstboot.service
After=systemd-homed-firstboot.service